# scripts/systemd/virtuoso-web.service
[Unit]
Description=Virtuoso Web Dashboard
Documentation=https://github.com/virtuoso/ccxt-trading
After=network-online.target virtuoso-trading.service
Requires=virtuoso-trading.service
Wants=network-online.target

[Service]
Type=simple
User=linuxuser
Group=linuxuser
WorkingDirectory=/home/linuxuser/trading/Virtuoso_ccxt

# Environment variables
EnvironmentFile=/home/linuxuser/trading/Virtuoso_ccxt/.env
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/linuxuser/trading/Virtuoso_ccxt"
Environment="PATH=/home/linuxuser/trading/Virtuoso_ccxt/venv311/bin:/usr/local/bin:/usr/bin:/bin"

# Startup
ExecStart=/home/linuxuser/trading/Virtuoso_ccxt/venv311/bin/python src/web_server.py

# Process management
Restart=always
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=20

# Kill all processes in the control group
KillMode=control-group
KillSignal=SIGTERM
SendSIGKILL=yes

# Resource limits (lighter than main engine)
LimitNOFILE=32768
MemoryMax=2G
CPUWeight=80

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=virtuoso-web

[Install]
WantedBy=multi-user.target