# scripts/systemd/virtuoso-log-cleanup.service
[Unit]
Description=Virtuoso CCXT log cleanup service
After=network.target

[Service]
Type=oneshot
User=linuxuser
Group=linuxuser
WorkingDirectory=/home/linuxuser/trading/Virtuoso_ccxt

# Environment variables
Environment="PYTHONPATH=/home/linuxuser/trading/Virtuoso_ccxt"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Main cleanup command
ExecStart=/usr/bin/python3 /home/linuxuser/trading/Virtuoso_ccxt/scripts/utilities/clean_logs.py --execute --logs-dir /home/linuxuser/trading/Virtuoso_ccxt/logs

# Log output
StandardOutput=append:/home/linuxuser/trading/Virtuoso_ccxt/logs/cleanup.log
StandardError=append:/home/linuxuser/trading/Virtuoso_ccxt/logs/cleanup_error.log

# Restart policy
Restart=no

# Resource limits
MemoryMax=256M
CPUQuota=50%
IOWeight=10

# Security settings
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/linuxuser/trading/Virtuoso_ccxt/logs /home/linuxuser/trading/Virtuoso_ccxt/archives

[Install]
WantedBy=multi-user.target