# scripts/systemd/virtuoso-trading.service
[Unit]
Description=Virtuoso Trading System
After=network-online.target memcached.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=linuxuser
Group=linuxuser
WorkingDirectory=/home/linuxuser/trading/Virtuoso_ccxt

# Environment
EnvironmentFile=/home/linuxuser/trading/Virtuoso_ccxt/.env
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/linuxuser/trading/Virtuoso_ccxt"

# Performance
CPUAffinity=0 1 2 3
Nice=-5
IOSchedulingClass=1
IOSchedulingPriority=4

# Execution
ExecStart=/home/linuxuser/trading/Virtuoso_ccxt/venv311/bin/python -u src/main.py

# Process management
Restart=always
RestartSec=10
KillMode=control-group

# Resource limits
LimitNOFILE=65536
MemoryMax=6G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=virtuoso

[Install]
WantedBy=multi-user.target