# scripts/systemd/virtuoso-health-monitor.service
[Unit]
Description=Virtuoso Trading Health Monitor
After=virtuoso-trading.service virtuoso-web.service
Wants=virtuoso-trading.service virtuoso-web.service

[Service]
Type=simple
User=linuxuser
Group=linuxuser
WorkingDirectory=/home/linuxuser

# Run health check script every 60 seconds
ExecStart=/bin/bash -c 'while true; do \
    if ! curl -s -f -m 5 http://localhost:8002/health > /dev/null 2>&1; then \
        echo "$(date): Web dashboard health check failed, restarting web service..."; \
        systemctl restart virtuoso-web.service; \
        sleep 30; \
    fi; \
    if ! curl -s -f -m 5 http://localhost:8001/api/monitoring/status > /dev/null 2>&1; then \
        echo "$(date): Monitoring API health check failed, restarting monitoring service..."; \
        systemctl restart virtuoso-monitoring-api.service; \
        sleep 30; \
    fi; \
    sleep 60; \
done'

Restart=always
RestartSec=30

# Kill all processes in the control group
KillMode=control-group
KillSignal=SIGTERM
SendSIGKILL=yes

# Resource limits
CPUQuota=5%
MemoryMax=100M

[Install]
WantedBy=multi-user.target