# scripts/systemd/virtuoso-optimized.service
[Unit]
Description=Virtuoso CCXT Trading System - Optimized for 4-core Singapore VPS
After=network.target memcached.service redis.service
Requires=memcached.service redis.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=linuxuser
Group=linuxuser
WorkingDirectory=/home/linuxuser/trading/Virtuoso_ccxt

# CPU Affinity - Pin to all 4 cores with optimal distribution
ExecStart=/usr/bin/taskset -c 0-3 /home/linuxuser/trading/Virtuoso_ccxt/venv/bin/python src/main.py
CPUAffinity=0 1 2 3

# Memory Limits - Allow 6GB for application processes
MemoryMax=6G
MemorySwapMax=2G

# Process Limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment Variables
Environment=PYTHONPATH=/home/linuxuser/trading/Virtuoso_ccxt
Environment=APP_PORT=8003
Environment=MONITORING_PORT=8001
Environment=CACHE_TYPE=memcached
Environment=MEMCACHED_HOST=localhost
Environment=MEMCACHED_PORT=11211
Environment=REDIS_HOST=localhost
Environment=REDIS_PORT=6379
Environment=LOG_LEVEL=INFO

# Process Management
Restart=always
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/linuxuser/trading/Virtuoso_ccxt
ReadWritePaths=/tmp
ReadWritePaths=/var/log

# I/O Priority for SSD optimization
IOSchedulingClass=1
IOSchedulingPriority=4

# CPU Scheduling
CPUSchedulingPolicy=0
CPUSchedulingPriority=0
Nice=-5

# Watchdog
WatchdogSec=60

[Install]
WantedBy=multi-user.target