# scripts/systemd/virtuoso-log-cleanup.timer
[Unit]
Description=Daily Virtuoso CCXT log cleanup timer
Requires=virtuoso-log-cleanup.service

[Timer]
# Run daily at 3:00 AM
OnCalendar=daily
OnCalendar=*-*-* 03:00:00

# Run weekly aggressive cleanup on Sunday at 2:00 AM
OnCalendar=weekly
OnCalendar=Sun *-*-* 02:00:00

# If the system was off when the timer should have run, run it when it starts
Persistent=true

# Randomize start time by up to 30 minutes to avoid load spikes
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target