# scripts/systemd/virtuoso-cache-warmer.service
[Unit]
Description=Virtuoso Cache Warmer Service
After=network.target virtuoso-trading.service
Wants=virtuoso-trading.service

[Service]
Type=simple
User=linuxuser
Group=linuxuser
WorkingDirectory=/home/linuxuser/trading/Virtuoso_ccxt
Environment="PYTHONPATH=/home/linuxuser/trading/Virtuoso_ccxt"
ExecStart=/home/linuxuser/trading/Virtuoso_ccxt/venv311/bin/python src/monitoring/cache_warmer.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
