# scripts/supervisor_config.ini
[program:virtuoso-main]
command=/home/linuxuser/trading/Virtuoso_ccxt/venv311/bin/python -u src/main.py
directory=/home/linuxuser/trading/Virtuoso_ccxt
user=linuxuser
autostart=true
autorestart=true
stderr_logfile=/home/linuxuser/trading/Virtuoso_ccxt/logs/main_error.log
stdout_logfile=/home/linuxuser/trading/Virtuoso_ccxt/logs/main.log
environment=PATH="/home/linuxuser/trading/Virtuoso_ccxt/venv311/bin:%(ENV_PATH)s",PYTHONPATH="/home/linuxuser/trading/Virtuoso_ccxt"

[program:virtuoso-web]
command=/home/linuxuser/trading/Virtuoso_ccxt/venv311/bin/uvicorn src.web_server:app --host 0.0.0.0 --port 8001
directory=/home/linuxuser/trading/Virtuoso_ccxt
user=linuxuser
autostart=true
autorestart=true
stderr_logfile=/home/linuxuser/trading/Virtuoso_ccxt/logs/web_error.log
stdout_logfile=/home/linuxuser/trading/Virtuoso_ccxt/logs/web.log
environment=PATH="/home/linuxuser/trading/Virtuoso_ccxt/venv311/bin:%(ENV_PATH)s",PYTHONPATH="/home/linuxuser/trading/Virtuoso_ccxt"

[group:virtuoso]
programs=virtuoso-main,virtuoso-web
priority=999