# scripts/deployment/virtuoso-fixed.service
[Unit]
Description=Virtuoso Trading System
After=network.target
Conflicts=virtuoso.service

[Service]
Type=simple
User=linuxuser
WorkingDirectory=/home/linuxuser/trading/Virtuoso_ccxt
Environment="PATH=/home/linuxuser/trading/Virtuoso_ccxt/venv311/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Use a lock file to prevent multiple instances
ExecStartPre=/bin/bash -c 'if [ -f /tmp/virtuoso.lock ]; then PID=$(cat /tmp/virtuoso.lock); if ps -p $PID > /dev/null; then echo "Already running with PID $PID"; exit 1; fi; fi'
ExecStartPre=/bin/bash -c 'echo $$ > /tmp/virtuoso.lock'
ExecStartPre=/bin/bash -c 'sudo lsof -ti:8003 | xargs -r sudo kill -9 || true'
ExecStartPre=/bin/sleep 2

ExecStart=/home/linuxuser/trading/Virtuoso_ccxt/venv311/bin/python -u src/main.py

ExecStopPost=/bin/rm -f /tmp/virtuoso.lock

# Restart configuration
Restart=on-failure
RestartSec=30
StartLimitBurst=3
StartLimitInterval=300

# Resource limits
MemoryMax=2G
CPUQuota=80%

# Logging
StandardOutput=journal+console
StandardError=journal+console
SyslogIdentifier=virtuoso

# Timeout for starting/stopping
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target