# scripts/deployment/virtuoso.service
[Unit]
Description=Virtuoso Trading System
After=network.target

[Service]
Type=simple
User=linuxuser
WorkingDirectory=/home/linuxuser/trading/Virtuoso_ccxt
Environment="PATH=/home/linuxuser/trading/Virtuoso_ccxt/venv311/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStartPre=/bin/bash -c 'sudo lsof -ti:8003 | xargs -r sudo kill -9 || true'
ExecStart=/home/linuxuser/trading/Virtuoso_ccxt/venv311/bin/python src/main.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/linuxuser/trading/Virtuoso_ccxt/app.log
StandardError=append:/home/linuxuser/trading/Virtuoso_ccxt/app.log

[Install]
WantedBy=multi-user.target